# scraper/interpol/crontab
*/5  * * * * root scrapy crawl wanted >> /var/log/cron.log 2>&1
